# henkaru-wiki-1.0/debian/cron.d
*/5 * * * * root /opt/wiki/bin/trafreports.py 2>&1 >/dev/null
